# lab/wordpress/run-mysql.ini
[program:mysql]
command=/usr/bin/mysqld --datadir=/var/lib/mysql
user=mysql
autostart=true
autorestart=true
