# lab/wordpress/run-apache.ini
[program:apache]
command=/usr/sbin/httpd -D FOREGROUND
user=root
autostart=true
autorestart=true

